# caddy-proxy/root/etc/services.d/docker-gen/finish
#!/usr/bin/with-contenv bash

s6-svscanctl -t /s6
