# homeassistant/root/etc/services.d/homeassistant/run
#!/usr/bin/with-contenv sh

python -m homeassistant --config /config
