# caddy-proxy/root/etc/services.d/caddy/run
#!/usr/bin/with-contenv bash

exec /usr/bin/caddy -conf /etc/Caddyfile $ENV_CADDY_OPTIONS
