# caddy-proxy/etc/supervisor.d/caddy.conf
[program:caddy]
command=caddy -conf /etc/Caddyfile %(ENV_CADDY_OPTIONS)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
